# projects/phoenix-agent/deployments/systemd/phoenix-agent.service
[Unit]
Description=Phoenix Agent - Observability Cost Optimization Agent
Documentation=https://github.com/phoenix/platform
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/phoenix-agent
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=phoenix-agent

# Environment configuration
Environment="PHOENIX_API_URL=http://phoenix-api:8080"
Environment="POLL_INTERVAL=15s"
Environment="LOG_LEVEL=info"
Environment="CONFIG_DIR=/etc/phoenix-agent"
EnvironmentFile=-/etc/phoenix-agent/environment

# Security settings (adjust based on requirements)
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/etc/phoenix-agent /var/log/phoenix-agent /tmp

# Resource limits
LimitNOFILE=65535
LimitNPROC=4096

# Kill settings
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target